# auth/auth.service
[Service]
ExecStart=/usr/sbin/auth
Environment=ETC=/etc/